# roles/prometheus/templates/node_exporter.service.j2
[Unit]
Description=Node Exporter
After=network.target

[Service]
User=root
Group=root
Type=simple
ExecStart=/usr/bin/prometheus-node-exporter

[Install]
WantedBy=multi-user.target